# ops/saturn/systemd/cicero-openclaw-gateway.service
[Unit]
Description=Cicero: OpenClaw Gateway (loopback only)
After=network.target

[Service]
Type=simple

# Assumes the repo is checked out at: %h/cicero
WorkingDirectory=%h/cicero

# IMPORTANT:
# Replace `<GATEWAY_SUBCOMMAND>` and `<FLAGS_FOR_HOST_AND_PORT>` with the real OpenClaw gateway start command.
#
# Discover available commands:
#   %h/cicero/bin/openclaw --help
#   %h/cicero/bin/openclaw <subcommand> --help
#
# Ensure it binds only to loopback and uses the intended port:
#   host: 127.0.0.1
#   port: 18789
ExecStart=%h/cicero/bin/openclaw <GATEWAY_SUBCOMMAND> <FLAGS_FOR_HOST_AND_PORT>

Restart=always
RestartSec=2

[Install]
WantedBy=default.target
